# deploy/exam-backend.service
[Unit]
Description=Exam System Backend Service
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/exam-backend

# JVM 参数优化
Environment="JAVA_OPTS=-Xms512m -Xmx1024m -XX:+UseG1GC -XX:MaxGCPauseMillis=200"

# 启动命令
ExecStart=/usr/bin/java $JAVA_OPTS -jar /opt/exam-backend/springboot-0.0.1-SNAPSHOT.jar

# 重启策略
Restart=on-failure
RestartSec=10

# 日志配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=exam-backend

# 安全配置
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
